# supervisord.conf
[supervisord]
nodaemon=true

[program:keep_running]
command=/bin/bash -c "while true; do sleep 1000; done"
autostart=true
autorestart=true

[program:zsh]
command=zsh
autostart=true
autorestart=true
